# systemd/rps-backup.service
[Unit]
Description=RPS Backup Service
Documentation=file:///home/paul/src/rps/docs/BACKUP_GUIDE.md
After=network.target

[Service]
Type=oneshot
User=paul
Group=paul
WorkingDirectory=/home/paul/src/rps

# Run the backup script
ExecStart=/home/paul/src/rps/bin/backup --keep 14

# Send notification on failure (requires mail/sendmail configured)
# ExecStartPost=-/bin/sh -c 'if [ "$EXIT_STATUS" != "0" ]; then echo "RPS backup failed" | mail -s "RPS Backup Failed" admin@example.com; fi'

# Environment
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/paul/src/rps/backups /home/paul/src/rps/logs /tmp

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rps-backup

[Install]
WantedBy=multi-user.target
